# _ kiwami/tools/daily-log/com.user.daily-log-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.daily-log-server</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/fukuokase/bin/daily-log/server/daily-log-server.py</string>
  </array>

  <key>KeepAlive</key>
  <true/>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/daily-log-server.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/daily-log-server-error.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>LANG</key>
    <string>ja_JP.UTF-8</string>
    <key>DAILY_LOG_HOST</key>
    <string>127.0.0.1</string>
    <key>DAILY_LOG_PORT</key>
    <string>8765</string>
  </dict>
</dict>
</plist>
